# setup/debian/speedtest_rs.socket
[Unit]
Description=librespeed rust backend (socket)
PartOf=speedtest_rs.service

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target